# Documentation/NodePraxisbuchBeispiele/Kapitel21/Listing21.3
#! /bin/sh
### BEGIN INIT INFO
# Provides:          Webserver in Node.js
# Required-Start:
# Required-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: HTTP-Server implementierung in Node.js
# Description:       Stellt einen Webserver auf Port 80 zur Verfügung
### END INIT INFO
# Author: Name <root@example.com>

# Aktionen
case "$1" in
    start)
        /srv/node/server.js start
        ;;
    stop)
        /srv/node/server.js stop
        ;;
    restart)
        /srv/node/server.js restart
        ;;
esac

exit 0
